# 6.apache2/update2.service
[Unit]
Description=Update2

[Service]
ExecStart=/bin/bash /var/www/html/update2.sh
Restart=always

[Install]
WantedBy=multi-user.target
